# launchd/local.ssh.socks.bc2.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>local.ssh.socks.bc2</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>/Users/takumi.aoki/commands/scripts/commands-45454-sock5-proxy-ag-claude.sh</string>
    </array>

    <!-- Restart automatically if the tunnel dies -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Throttle restart: wait 10 s before respawning after a crash -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Discard stdout/stderr; the script itself emits nothing useful -->
    <key>StandardOutPath</key>
    <string>/dev/null</string>
    <key>StandardErrorPath</key>
    <string>/dev/null</string>
  </dict>
</plist>
